# system/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
pidfile=/tmp/supervisord.pid
minfds=1024
minprocs=200
nodaemon=true

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:nginx]
command=/srv/system/run_nginx.sh
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:app_server]
command=/srv/system/run_app_server.sh
startsecs=5
environment=PYTHONPATH=/srv
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:cron]
command = /bin/bash -c "declare -p | grep -Ev '^declare -[[:alpha:]]*r' > /run/supervisord.env && /usr/sbin/cron -f -L 15"
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
autorestart = true
priority = 20

[program:webpack]
command=/srv/system/run_webpack_server.sh
startsecs=0
directory=/srv/ui
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:storybook]
command=/srv/system/run_storybook.sh
startsecs=0
directory=/srv/ui
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:watch_photos]
command=bash -c "nice -n 16 python /srv/photonix/manage.py watch_photos"
startsecs=10
environment=PYTHONPATH=/srv
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:raw_scheduler]
command=bash -c "sleep 5 && nice -n 17 python /srv/photonix/manage.py raw_scheduler"
startsecs=15
environment=PYTHONPATH=/srv
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:raw_processor]
command=bash -c "sleep 6 && nice -n 17 python /srv/photonix/manage.py raw_processor"
startsecs=16
environment=PYTHONPATH=/srv
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:thumbnail_scheduler]
command=bash -c "sleep 7 && nice -n 17 python /srv/photonix/manage.py thumbnail_processor"
startsecs=17
environment=PYTHONPATH=/srv
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:classification_scheduler]
command=bash -c "sleep 8 && nice -n 18 python /srv/photonix/manage.py classification_scheduler"
startsecs=18
environment=PYTHONPATH=/srv
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:classification_color_processor]
command=bash -c "sleep 9 && nice -n 19 python /srv/photonix/manage.py classification_color_processor"
startsecs=19
environment=PYTHONPATH=/srv
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:classification_location_processor]
command=bash -c "sleep 10 && nice -n 19 python /srv/photonix/manage.py classification_location_processor"
startsecs=20
environment=PYTHONPATH=/srv
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:classification_face_detection_processor]
command=bash -c "sleep 11 && nice -n 19 python /srv/photonix/manage.py classification_face_processor"
startsecs=21
environment=PYTHONPATH=/srv
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:classification_style_processor]
command=bash -c "sleep 12 && nice -n 19 python /srv/photonix/manage.py classification_style_processor"
startsecs=22
environment=PYTHONPATH=/srv
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:classification_object_processor]
command=bash -c "sleep 13 && nice -n 19 python /srv/photonix/manage.py classification_object_processor"
startsecs=23
environment=PYTHONPATH=/srv
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:classification_event_processor]
command=bash -c "sleep 14 && nice -n 19 python /srv/photonix/manage.py classification_event_processor"
startsecs=24
environment=PYTHONPATH=/srv
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0

[program:rescan_photos_periodically]
command=bash -c "sleep 60 && nice -n 16 python /srv/photonix/manage.py rescan_photos_periodically"
startsecs=70
environment=PYTHONPATH=/srv
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
stderr_logfile_maxbytes=0
stdout_logfile_maxbytes=0
