# system/cron.d/retrain_face_similarity_index
*/5 * * * * root exec /bin/bash -c ". /run/supervisord.env; python /srv/photonix/manage.py retrain_face_similarity_index"
